# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:queue-worker]
command=php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=60
directory=/var/www
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/supervisor/queue-worker.log
stderr_logfile=/var/log/supervisor/queue-worker-error.log
stopwaitsecs=70

[program:schedule-runner]
command=sh -c "while true; do php /var/www/artisan schedule:run >> /dev/null 2>&1; sleep 60; done"
directory=/var/www
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/supervisor/schedule-runner.log
stderr_logfile=/var/log/supervisor/schedule-runner-error.log

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true